# root/etc/services.d/bot/run
#!/usr/bin/with-contenv bash

exec \
	s6-setuidgid bot \
		python /app/bot.py
